# projector_supervisor.conf
[program:projector_experiments]
command=python app.py
user=root
directory=/opt/projector_experiments
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=INT
